# systemd/system/backup.service
[Unit]
Description=Backup system service

[Service]
Type=simple
ExecStart=/home/kyle/.bin/backup

